# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/supervisord.log

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:php-fpm]
command=php-fpm -F -R
